# services/redis-charlie-node/run
#!/bin/sh

exec redis-server --port 6384
