# services/.logger/run
#!/bin/sh

out="$LOGS_DIR/$NAME"

mkdir -p "$out"

exec svlogd "$out"
