# services/.redis/run
#!/bin/sh

. config.sh

exec redis-server --port $PORT
